# packaging/systemd/keyflow.service
[Unit]
Description=KeyFlow — Non-paste password input assistant
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/keyflow run
Restart=on-failure
RestartSec=3
PassEnvironment=DISPLAY XAUTHORITY BW_SESSION BW_PASSWORD
EnvironmentFile=-%h/.config/keyflow/env

[Install]
WantedBy=graphical-session.target
